# systemd/thephone.service
[Unit]
Description=The Phone
After=multi-user.target

[Service]
Type=simple
ExecStart=/home/op/thephone
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
